# src/init/sv/sshd/run
#!/bin/sh
# Dropbear SSH server (disk-boot / runit path).
# -F: foreground (runit supervises it).  -E: log to stderr (svlogd captures it).
# -R: generate host keys on demand, stored in /etc/dropbear.
mkdir -p /etc/dropbear
exec /usr/sbin/dropbear -F -E -R -p 22
