# src/init/sv/getty-tty1/run
#!/bin/sh
exec /sbin/getty -L tty1 115200 vt100
