# src/init/sv/crond/run
#!/bin/sh
exec /usr/sbin/crond -n -c /var/spool/cron/crontabs
